# .provision/flask_project.conf
[program:flask_project]
command = gunicorn app:app -b localhost:8080 --reload
directory = /vagrant
user = ubuntu
